# mint-backend/systemd/corecats-mint-backend.service.example
[Unit]
Description=Core Cats mint backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/core-cats/mint-backend
EnvironmentFile=/etc/corecats-mint-backend.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m corecats_mint_backend.server
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
